# helpers/debianpackaging/package/DEBIAN/postinst
if [ ! -f /var/lib/kaithem/users/data/users.json ]
then
    adduser --system --no-create-home kaithem
    usermod -a -G serial kaithem
    usermod -a -G dialout kaithem
    usermod -a -G audio kaithem
    usermod -a -G gpio kaithem
    usermod -a -G i2c kaithem
    usermod -a -G video kaithem
    usermod -a -G spi kaithem
    groupadd -r kaithem
    usermod -a -G kaithem kaithem
    mkdir /var/log/kaithem
    chmod 700 /var/log/kaithem
    chown kaithem /var/log/kaithem
    mkdir /var/lib/kaithem
    chown kaithem /var/lib/kaithem
    chmod 700 /var/lib/kaithem/
    mkdir /var/lib/kaithem/ssl
    chown kaithem /var/lib/kaithem/ssl
    chmod 700 /var/lib/kaithem/ssl
    #Look at this race condition here! Someone could steal that private key before we chmod it!
    openssl req -new -newkey rsa:2048 -days 365 -nodes -x509 -keyout /var/lib/kaithem/ssl/certificate.key  -out  /var/lib/kaithem/ssl/certificate.cert -subj '/CN=localhost'
    chmod 600 /var/lib/kaithem/ssl/certificate.key
    chown kaithem /var/lib/kaithem/ssl/certificate.key
    systemctl enable kaithem
    kaithem --initialpackagesetup x
    systemctl start kaithem
fi